# clash/clash.service
[Unit]
Description=Clash Daemon

[Service]
ExecStart=/data/clash/clash -d /data/clash/
Restart=on-failure

[Install]
WantedBy=multi-user.target
